# subsonic-installer-debian/src/DEBIAN/prerm
#! /bin/sh

set -e

# Stop Subsonic service.
[ -e /etc/init.d/subsonic ] && invoke-rc.d subsonic stop
