# subsonic-installer-debian/src/main/resources/DEBIAN/postinst
#! /bin/sh

set -e

ln -s /usr/share/subsonic/subsonic.sh /usr/bin/subsonic

chmod oug+rwx /var/subsonic
rm -rf /var/subsonic/jetty

update-rc.d subsonic defaults 99
invoke-rc.d subsonic restart
