# subsonic-installer-debian/src/main/resources/DEBIAN/postrm
#! /bin/sh

set -e

rm -f /usr/bin/subsonic

[ -e /etc/init.d/subsonic ] && invoke-rc.d subsonic stop
update-rc.d -f subsonic remove