# subsonic-installer-debian/src/main/resources/DEBIAN/postrm
#! /bin/sh

set -e

if [ -h /usr/bin/subsonic ] ; then
  rm /usr/bin/subsonic
fi

invoke-rc.d subsonic stop
update-rc.d -f subsonic remove